# Web/basic/deployment/bin/app.service
[Unit]
Description=Golang Application
Requires=network-online.target
After=network-online.target

[Service]
User=app
Group=app
Restart=on-failure
ExecStart=/usr/bin/app
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target